# docker/web/etc/cron.hourly/process-activity-stats
#!/bin/bash

cd /app/backend || exit
php yii user/process-activity-stats --interactive=0
